# rootfs/etc/supervisord/conf.d/cron.conf
[program:cron]
command = cron -f
priority=100
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
