# rootfs/etc/supervisord/conf.d/xshokfix.conf
[program:xshok-fix]
command = /xshok-fix.sh
startsecs = 0
autorestart = false
startretries = 1
